# printer-manager/src/files/export.cron
0 0 1 * * root .utils.sh --set-monthly-cote
